# irdeluxe.service
[Unit]
Description=IR-REST gateway

[Service]
Type=simple
Restart=always
ExecStart=/root/ircontroller/start.sh

[Install]
WantedBy=multi-user.target
